# root/etc/s6-overlay/s6-rc.d/svc-pure-ftpd/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "[pure-ftpd] Starting with flags:$PUREFTPD_FLAGS"

# shellcheck disable=SC2086
exec pure-ftpd ${PUREFTPD_FLAGS}
