# root/etc/s6-overlay/s6-rc.d/init-pure-ftpd-config/up
/etc/s6-overlay/s6-rc.d/init-pure-ftpd-config/run
