# root/etc/s6-overlay/s6-rc.d/svc-syslog/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec syslogd -ntSD -O - -s 0
